# hpc/crontab/crontab_task

# every monday
0 0 * * 1 bash /path/to/weekGetUsage.sh
